# deploy/gunicorn.service
[Unit]
Description=Slgi Django (gunicorn)
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=wlgusqkr22
Group=wlgusqkr22
WorkingDirectory=/home/wlgusqkr22/slgi/backend
Environment="PATH=/home/wlgusqkr22/.local/bin:/usr/local/bin:/usr/bin:/bin"
Environment="DJANGO_SETTINGS_MODULE=config.settings.production"
EnvironmentFile=/home/wlgusqkr22/slgi/backend/.env
ExecStart=/home/wlgusqkr22/.local/bin/uv run gunicorn config.wsgi:application \
    --bind 127.0.0.1:8000 \
    --workers 3 \
    --timeout 60 \
    --access-logfile - \
    --error-logfile -
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
